# provisioning/roles/install-slack-bots/templates/slackin.service.j2
[Unit]
Description=slackin

[Service]
ExecStart=/usr/bin/slackin opendatakit {{ lookup('file', 'secrets/slackin.token') }}
Restart=on-failure

[Install]
WantedBy=multi-user.target